# seuds-discord-bot.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
PWD=/home/bot/code
# Checks the bot every 10 minutes
*/10 * * * * bot ./launch.sh start